# samples/revor.supervisor.conf
[program:revor]
command=/opt/revor/current/run.sh
user=root
autostart=true
autorestart=true