# samples/revor.supervisor.conf
[program:revor]
command=/opt/revor/current/revor -importPath revor -srcPath /opt/revor/current/src -runMode prod
directory=/opt/revor/current
user=root
autostart=true
autorestart=true