# samples/revor.supervisor.conf
[program:revor]
command=/opt/revor/cirrent/revor
user=root
autostart=true
autorestart=true